# task1/unit-files/check-status.service
[Unit]
Description=Check status service

[Service]
Type=oneshot
ExecStart=/opt/webapp/script.sh
